# netguard_scan.timer
[Unit]
Description=NetGuard Network Scanner Timer

[Timer]
OnBootSec=60
OnUnitActiveSec=300s

[Install]
WantedBy=timers.target
