# netguard_scan.service
[Unit]
Description=NetGuard Network Scanner Service
After=network.target postgresql.service

[Service]
Type=oneshot
User=netguard
WorkingDirectory=/opt/netguard
ExecStart=/opt/netguard/venv/bin/python net_scan.py
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
