# netureon_scan.timer
[Unit]
Description=Run Netureon network scan every 5 minutes

[Timer]
# Wait 1 minute after boot before first scan
OnBootSec=1min
# Run the scan every 5 minutes
OnUnitActiveSec=5min
# Ensure the timer triggers with an accuracy of 1 second
AccuracySec=1s
# Ensure timer is persistent across reboots
Persistent=true

[Install]
WantedBy=timers.target
