# netureon.service
[Unit]
Description=Netureon Network Security Monitor
After=network.target postgresql.service
Conflicts=netureon_scan.service

[Service]
Type=notify
User=orangepi
Group=orangepi
WorkingDirectory=/home/orangepi/Netureon
Environment=PATH=/home/orangepi/Netureon/.venv/bin:$PATH
ExecStart=/home/orangepi/Netureon/.venv/bin/python3 /home/orangepi/Netureon/net_scan.py
ExecStop=/bin/kill -TERM $MAINPID
RemainAfterExit=no
Restart=always
RestartSec=30
TimeoutStartSec=120
TimeoutStopSec=30
WatchdogSec=60
NotifyAccess=all

# Set resource limits
CPUQuota=50%
MemoryMax=512M

[Install]
WantedBy=multi-user.target
