# netureon.service
[Unit]
Description=Netureon Network Monitor
After=network.target postgresql.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h/Netureon
Environment=PATH=%h/Netureon/.venv/bin:$PATH
ExecStart=%h/Netureon/.venv/bin/python main.py
Restart=always
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
