# netureon-alerts.service
[Unit]
Description=Netureon Network Security Monitor
After=network.target postgresql.service

[Service]
Type=simple
User=orangepi
Group=orangepi
WorkingDirectory=/home/orangepi/Netureon
Environment=PYTHONPATH=/home/orangepi/Netureon
ExecStart=/home/orangepi/Netureon/.venv/bin/python3 alert_daemon.py
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
