# netureon_web.service
[Unit]
Description=Netureon Web Interface
After=network.target postgresql.service

[Service]
Type=simple
User=orangepi
WorkingDirectory=%h/Netureon
Environment=PYTHONPATH=%h/Netureon
ExecStart=%h/Netureon/.venv/bin/python webui/app.py
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
