# netureon_scan.service
[Unit]
Description=Netureon Network Scanner Service
After=network.target postgresql.service
Wants=netureon_scan.timer

[Service]
Type=oneshot
User=orangepi
WorkingDirectory=/home/orangepi/Netureon
Environment=PATH=/home/orangepi/Netureon/.venv/bin:$PATH
ExecStart=/home/orangepi/Netureon/.venv/bin/python net_scan.py
RemainAfterExit=no
# Set resource limits to prevent runaway scans
CPUQuota=50%
MemoryMax=512M
# Add timeout to prevent stuck scans
TimeoutSec=300
# Basic error handling
Restart=on-failure
RestartSec=30s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
